# configs/supervisor/supervisord-selenium.conf
[program:selenium]
command=/start-selenium.sh

